# deploy/systemd/claw-plaid-ledger-sync.timer
# claw-plaid-ledger-sync.timer — primary sync trigger for claw-plaid-ledger-sync.service
#
# This timer is the PRIMARY sync mechanism for claw-plaid-ledger.  It replaces
# the former webhook-based sync path (disabled by default as of M25).
# See RUNBOOK.md Section 12.3 for details.
#
# Activation:
#   systemctl enable --now claw-plaid-ledger-sync.timer
#
# Status check:
#   systemctl list-timers claw-plaid-ledger-sync.timer

[Unit]
Description=Run claw-plaid-ledger sync four times daily
Documentation=https://github.com/aronfoster/claw-plaid-ledger

[Timer]
# Default: four syncs per day (midnight, 06:00, noon, 18:00).
# To sync hourly instead, create a drop-in override:
#   systemctl edit claw-plaid-ledger-sync.timer
# and add:
#   [Timer]
#   OnCalendar=
#   OnCalendar=hourly
# The first empty OnCalendar= clears the default before setting hourly.
OnCalendar=*-*-* 00,06,12,18:00:00

# Persistent=true means the timer will catch up a missed run (e.g. after
# the host was powered off) as soon as it boots rather than waiting for the
# next scheduled slot.  This ensures syncs are not silently skipped.
Persistent=true

# Optional: add a small random delay to spread load if multiple units fire
# at the same time.  Disabled by default.
# RandomizedDelaySec=120

# The service unit to activate when the timer fires.
Unit=claw-plaid-ledger-sync.service

[Install]
# Timer units are enabled directly; they pull in the service unit themselves.
WantedBy=timers.target
