# deploy/systemd/claw-plaid-ledger-sync.service
# claw-plaid-ledger-sync.service — oneshot unit for "ledger sync --all"
#
# This unit is meant to be triggered by claw-plaid-ledger-sync.timer rather
# than started directly.  It runs a single synchronisation pass across all
# configured items and then exits.
#
# Installation:
#   cp deploy/systemd/claw-plaid-ledger-sync.service /etc/systemd/system/
#   cp deploy/systemd/claw-plaid-ledger-sync.timer   /etc/systemd/system/
#   systemctl daemon-reload
#   systemctl enable --now claw-plaid-ledger-sync.timer
#
# See RUNBOOK.md Section 12 for the trade-offs between this timer and the
# CLAW_SCHEDULED_SYNC_ENABLED built-in fallback.

[Unit]
Description=claw-plaid-ledger scheduled sync (ledger sync --all)
Documentation=https://github.com/aronfoster/claw-plaid-ledger
# Only run when the main server is active — avoids syncing while the DB
# is being initialised or the server is deliberately stopped.
# Remove this line if you run sync independently of "ledger serve".
After=network-online.target
Wants=network-online.target

[Service]
# Type=oneshot: systemd waits for the process to exit before marking the
# unit as "done", which is what we want for a scheduled batch job.
Type=oneshot

# --- User / Group -------------------------------------------------------
# Must match the owner of the config directory and access tokens.
User=alice
Group=alice

# --- Secrets and configuration ------------------------------------------
# Same EnvironmentFile used by claw-plaid-ledger.service.
# Adjust path to match your OS/user (see claw-plaid-ledger.service for notes).
EnvironmentFile=/home/alice/.config/claw-plaid-ledger/.env

# --- Executable ---------------------------------------------------------
# Adjust path to match your ledger installation (see claw-plaid-ledger.service).
ExecStart=/usr/local/bin/ledger sync --all

# --- Logging ------------------------------------------------------------
StandardOutput=journal
StandardError=journal
SyslogIdentifier=claw-plaid-ledger-sync

# --- Security hardening -------------------------------------------------
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/home/alice/.local/share/claw-plaid-ledger /home/alice/.local/state/claw-plaid-ledger

# Note: [Install] is intentionally absent — this unit is activated by its
# timer, not by systemctl enable directly.
